# packages/tailscale/tailscale.post-install
#!/bin/sh

/home/root/.vellum/bin/mount-rw

echo "installing tailscaled.service..."
cp /home/root/.vellum/share/tailscale/systemd/tailscaled.service /etc/systemd/system/tailscaled.service

echo "reloading systemd daemon"
systemctl daemon-reload

echo "enabling tailscaled service..."
systemctl enable tailscaled.service

/home/root/.vellum/bin/mount-restore

echo "restarting tailscaled service..."
systemctl restart tailscaled.service
