# packages/tailscale/tailscale.pre-deinstall
#!/bin/sh

/home/root/.vellum/bin/mount-rw

echo "stopping tailscaled service..."
systemctl disable --now tailscaled.service

echo "removing tailscaled.service..."
rm -f /etc/systemd/system/tailscaled.service

/home/root/.vellum/bin/mount-restore

echo "reloading systemd daemon"
systemctl daemon-reload

if [ "$VELLUM_PURGE" = "1" ]; then
    echo "deleting tailscale config files..."
    rm -rf /home/root/.local/var/lib/tailscale
    rm -rf /home/root/.config/tailscale
fi
