# packages/tailscale/tailscale.pre-deinstall
#!/bin/sh

echo "stopping tailscaled service..."
systemctl disable --now tailscaled.service

if [ "$VELLUM_PURGE" = "1" ]; then
    echo "deleting tailscale config files..."
    rm -rf /home/root/.local/var/lib/tailscale
    rm -rf /home/root/.config/tailscale
fi
