# packages/tripletap/tripletap.pre-deinstall
#!/bin/sh

# Stop and disable service before removal
systemctl disable xovi-tripletap --now 2>/dev/null || true

echo "tripletap service disabled"
