# packages/entware-rc/entware-rc.pre-deinstall
#!/bin/sh
/home/root/.vellum/bin/mount-rw
rcctl list | xargs -I {} systemctl disable --now entware-rc@{}
rm -f /etc/systemd/system/entware-rc@.service
systemctl daemon-reload
/home/root/.vellum/bin/mount-restore
